# root/etc/crontabs/abc
3  */8  * * *   /app/zap2xml.sh >> /config/log/zap2xml/zap2xml.log 2>&1
