# root/etc/crontabs/abc
3  */8  * * *   /app/zap2xml.sh >> /config/zap2xml.log 2>&1
